# caelestia/timers/battery-warning.timer
[Unit]
Description=Battery Warning Timer

[Timer]
OnBootSec=30sec
OnUnitActiveSec=30sec

[Install]
WantedBy=default.target
